# httpd/debian/symbiosis-httpd.cron.weekly
#!/bin/bash

# regenerate dhparams for httpd
symbiosis-generate-dhparams /var/lib/symbiosis/httpd
