# phpmyadmin/debian/postinst
#!/bin/bash

set -e

#
# Skip, if we are not in "configure" state
#
if [ "$1" != "configure" ]; then
        echo "I: Skipping configuration"
        exit 0
fi


#
#  Remove old files
#
for i in phpmyadmin.vhost.conf phpadmin.conf phpmyadmin.conf ; do
  if [ -e /etc/apache2/conf.d/$i ]; then
    rm -f /etc/apache2/conf.d/$i
  fi
done


#
#  Ensure that the PHPMyAdmin configuration file is included.
#
if [ -e /etc/phpmyadmin/apache.conf ]; then
    ln -sf /etc/phpmyadmin/apache.conf /etc/apache2/conf.d/phpmyadmin.conf
fi

#
#  Append Symbiosis configuration file requirement
#
if [ -e /etc/phpmyadmin/config.inc.php ]; then

    if ( ! grep -q "require('/etc/phpmyadmin/config.symbiosis.inc.php')" /etc/phpmyadmin/config.inc.php ); then
        cat >>/etc/phpmyadmin/config.inc.php <<EOF

/**
 * Read symbiosis config.
 */
if (is_readable('/etc/phpmyadmin/config.symbiosis.inc.php')) {
    require('/etc/phpmyadmin/config.symbiosis.inc.php');
}
EOF
    fi
fi

#
#  Reload Apache2 to ensure this change is picked up.
#
invoke-rc.d apache2 reload || true

#DEBHELPER#
exit 0
