# monit/system/symbiosis-monit-failure-email.service
[Unit]
Description=Notify by email about symbiosis-monit failure

[Service]
Type=simple
ExecStart=/usr/bin/symbiosis-monit-failure-email
