# monit/system/symbiosis-monit-failure-email.service
[Unit]
Description=Notify by email about symbiosis-monit failure

[Service]
Type=simple
ExecStart=/usr/sbin/symbiosis-monit-failure-email

