# crack/etc/cron.d/symbiosis-crack
#
#  Test the strength of user passwords.
#
#  Output will be sent to root by email.
#


# hourly check
1  * * * * root [ -x /usr/bin/test-symbiosis-passwords ] && /usr/bin/test-symbiosis-passwords --hourly
# weekly check
47 6 * * 7 root [ -x /usr/bin/test-symbiosis-passwords ] && /usr/bin/test-symbiosis-passwords --weekly


