# meta/update-motd.d/00-bytemark-symbiosis
#!/bin/sh

echo "\033[1;34;49m   ┌───────────┐  ┌────────────────────────┐ ┌────┐   ┌─┐   ┌───────┐┌──┐ ┌──┐
┌─┐│  ┌─┐  ├─┐ └┐┌┘ ┌┴──┐  ┌──┐  ┌────┐    │ │    │  ┌┘ └┐  │  ┌─┐  └┤  │┌┘ ┌┘
└─┘│  └─┘  └┐└┐ └┘ ┌┘   │  │  │  └───┐│  ├┐└┬┘┌┤  │ ┌┘┌─┐└┐ │  └─┘  ┌┤  └┘ ┌┘
┌─┐│  ┌──┐  │ └┐  ┌┘    │  │  │  ┌───┘│  ││   ││  │┌┘ └─┘ └┐│  ┌─┐  └┤  ┌─┐└─┐
└─┘│  └──┘ ┌┘  │  │     │  │  │  └────┘  │└┐ ┌┘│  ├┘  ┌─┐  └┤  │ └┐  │  │ └┐ └┐
   └───────┘   └──┘     └──┘  └──────────┘ └─┘ └──────┘ └──────┘  └─────┘  └──┘
\033[1;39m  ┌───                 o           o         \033[0;35m┌───
\033[1;39m  │               │                          \033[0;35m│    │              │         │
\033[1;39m  └──┐ │  │ ┌─┬─┐ ├──┐ │ ┌──┐ ┌─── │ ┌───    \033[0;35m└──┐ ├─── ┌─── ┌──┐ ├─── ┌─── ├──┐
\033[1;39m     │ │  │ │ │ │ │  │ │ │  │ └──┐ │ └──┐    \033[0;35m   │ │    │    ├──┘ │    │    │  │
\033[1;39m  ───┘ └──┤ │ │ │ └──┘ │ └──┘ ───┘ │ ───┘    \033[0;35m───┘ └─── │    └─── └─── └─── │  │
\033[1;39m       ───┘\033[0m
  For documentation please visit \033[4;34mhttps://symbiosis.bytemark.co.uk/stretch/docs\033[0m
"
