# meta/update-motd.d/00-bytemark-symbiosis
#!/bin/sh

cat <<EOF
                                          ,o'o'o.
                                      _,-o,-'''-o:o.
  _ _____   _______ ___ __  __   _  -'o,o_  __      o
 (_) _ ) \\ / /_   _| __|  \\/  | /_\\ | _ \\ |/ /
  _| _ \\\\ V /  | | | _|| |\\/| |/ _ \\|   /   <
 (_)___/ |_|   |_| |___|_| _|_|_/ \\_\\_|_\\_|\\_\\
   / __/__ __ __ _  / /   (_)___   ___  (_)___
  _\\ \\ / // //  ' \\/ _ \\ / // _ \\ (_-< / /(_-<            ,o..           o
 /___/ \\_, //_/_/_/_.__//_/ \\___//___//_//___/     ,,,o-o'o.--o:o:o,,..:o
     /___/                         \`--o.o\`o---o'o'o,o,-'''        o'o'o
                                         \`\`-o-o''''
  For documentation please visit
    https://symbiosis.bytemark.co.uk/stretch/docs

EOF

