# common/system/symbiosis-skel.service
[Unit]
Description=Service for symbiosis-skel. Activated by symbiosis-skel.path

[Service]
ExecStart=/usr/bin/symbiosis-skel 
User=admin
Group=admin
