# common/system/symbiosis-skel.path
[Unit]
Description=Watcher for symbiosis-skel automated domain population

[Path]
PathExistsGlob=/srv/*

[Install]
WantedBy=multi-user.target
