# monit/cron.d/bytemark-vhost-monit
0-59/2 * * * *   root  /usr/sbin/bytemark-vhost-monit -m -d /etc/bytemark-vhost/monit.d
