# monit/cron.d/symbiosis-monit
0-59/2 * * * *   root  [ -x /usr/sbin/symbiosis-monit ] &&  /usr/sbin/symbiosis-monit -m -d /etc/symbiosis/monit.d
