# meta/debian/postinst
#!/bin/bash

set -e

#
# Skip, if we are not in "configure" state
#
if [ "$1" != "configure" ]; then
    echo "I: Skipping configuration"
    exit 0
fi


#
#  If there is volatile then comment it out.
#
if ( grep -q volatile /etc/apt/sources.list 2>/dev/null >/dev/null ); then
  sed -i -e 's/^[^#].*volatile.*$/# Removed since volatile no longer exists for squeeze\n# \0/g' /etc/apt/sources.list
fi


#
#  Update /etc/motd.
#
cat >> /etc/motd.tail <<EOF

To access the Bytemark Symbiosis documentation, please visit:

  * http://symbiosis.bytemark.co.uk/

EOF


#
#  Remove our obsolete vhost init script.
#
if [ -e /etc/init.d/bytemark-vhost ]; then
    update-rc.d -f bytemark-vhost remove
    rm /etc/init.d/bytemark-vhost
fi


#DEBHELPER#
exit 0
