# email/debian/postinst
#!/bin/sh

set -e

#
# Skip, if we are not in "configure" state
#
if [ "$1" != "configure" ]; then
        echo "I: Skipping configuration"
        exit 0
fi

#
#  Previous mistake
#
if [ -e /etc/cron.d/exim_rewrite_scan ]; then
    rm -rf /etc/cron.d/exim_rewrite_scan
fi
if [ -e /etc/cron.d/bytemark-vhost-email ]; then
    rm -rf /etc/cron.d/bytemark-vhost-email
fi

#
#  Add the user
#
adduser clamav Debian-exim 

#
#  Enable the spamassassin module
#
sed -i -e 's/^\(ENABLED=\)0/\11/' /etc/default/spamassassin

#
# Enable spamassassin daily updates
#
sed -i -e 's/^\(CRON=\)0/\11/' /etc/default/spamassassin

#
#  Make sure ownership is correct on our file
#
touch /etc/exim4/rewrites
chown Debian-exim /etc/exim4/rewrites

#
# Copy any existing files in /etc/exim4/bytemark-vhost.d, backing up anything
# in the way.
#
if [ -d /etc/exim4/bytemark-vhost.d ] ; then
  cp -a -b -S .dpkg-new /etc/exim4/bytemark-vhost.d/* /etc/exim4/symbiosis.d
  rm -rf /etc/exim4/bytemark-vhost.d
fi


#
#  Rebuild exim4
#
if [ -e /etc/exim4/Makefile ]; then
    cd /etc/exim4 && make
fi


#
#  Rebuild dovecot
#
if [ -e /etc/dovecot/Makefile ] ; then
    cd /etc/dovecot && make
fi

for i in spamassassin clamav-daemon exim4 dovecot; do

    #
    # Restart
    #
    if which invoke-rc.d >/dev/null 2>&1; then
        invoke-rc.d $i restart
    else
        /etc/init.d/$i restart
    fi
done

#DEBHELPER#
exit 0
