# email/debian/symbiosis-email.symbiosis-email-poppassd.init
#!/bin/sh
# kFreeBSD do not accept scripts as interpreters, using #!/bin/sh and sourcing.
if [ true != "$INIT_D_SCRIPT_SOURCED" ] ; then
    set "$0" "$@"; INIT_D_SCRIPT_SOURCED=true . /lib/init/init-d-script
fi
### BEGIN INIT INFO
# Provides:          symbiosis-email-poppassd
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Symbiosis: email poppass daemon
# Description:       Service used to change email passwords using ye olde poppass protocol
### END INIT INFO

# Author: Patrick J Cherry <patrick@bytemark.co.uk>

DESC="Symbiosis: email dict proxy"
DAEMON=/usr/sbin/symbiosis-email-poppassd
NAME="ruby"


