# email/debian/symbiosis-email.cron.weekly
#!/bin/bash -eu

verbose=""
if [ "$*" = "--verbose" ] || [ "$*" = "-v" ]; then
  verbose="-v "
fi

for i in /var/spool/exim4/gnutls-params-[0-9]* ; do
  rm -f "$i"
done
symbiosis-generate-dhparams $verbose/var/spool/exim4/gnutls-params-2048 Debian-exim

