# email/debian/symbiosis-email.cron.weekly
#!/bin/bash -eu

symbiosis-generate-dhparams /var/spool/exim4

for i in /var/spool/exim4/gnutls-params-[0-9]* ; do
  rm -f "$i"
done

mv /var/spool/exim4/dhparams /var/spool/exim4/gnutls-params-2048
