# webmail/debian/prerm
#!/bin/bash

set -e


#
# Remove our local alias
#
rm -f /etc/apache2/conf.d/squirrelmail.local.conf


#
#  Restart Apache2.
#
invoke-rc.d apache2 force-reload


#DEBHELPER#
exit 0
