# webmail/debian/symbiosis-webmail-roundcube.prerm
#!/bin/bash

set -e

#
# Remove our alternative
#
update-alternatives --remove symbiosis-webmail /var/lib/roundcube/

#
# Restart apache2 if apache2ctl exists and is executable.
#
if [ -x /usr/sbin/apache2ctl ] ; then
  #
  #  Restart Apache2.
  #
  invoke-rc.d apache2 force-reload || true
fi

#DEBHELPER#
exit 0
