# monit/debian/symbiosis-monit.cron.d
0-59/2 * * * *   root [ -x /usr/sbin/symbiosis-monit ] && /usr/sbin/symbiosis-monit -t email /etc/symbiosis/monit.d
