# monit/debian/symbiosis-monit.service
[Unit]
Description=Symbiosis monitoring

[Service]
Type=oneshot
ExecStart=/usr/sbin/symbiosis-monit
