# monit/debian/symbiosis-monit.service
[Unit]
Description=Symbiosis monitor

[Service]
Type=simple
ExecStart=/usr/sbin/symbiosis-monit -t email /etc/symbiosis/monit.d
