# xmpp/debian/symbiosis-xmpp.cron.weekly
#!/bin/bash -eu

verbose=""
if [ "$*" = "--verbose" ] || [ "$*" = "-v" ]; then
  verbose="-v "
fi

symbiosis-generate-dhparams $verbose/var/spool/prosody/dhparams prosody
