# xmpp/debian/symbiosis-xmpp.cron.weekly
#!/bin/bash -eu

symbiosis-generate-dhparams /var/spool/prosody prosody
