# firewall/cron.d/firewall-blocker
# Crontab snippet which will invoke our firewall-based blocker to firewall
# away people who conduct dictionary attacks.

*/5 * * * * root [ -x /usr/bin/firewall-blacklist ] && /usr/bin/firewall-blacklist

