# ftpd/debian/postinst
#!/bin/bash

set -e

#
# Skip, if we are not in "configure" state
#
if [ "$1" != "configure" ]; then
        echo "I: Skipping configuration"
        exit 0
fi

#
#  Remove existing PAMAuthentication setup,
#
find /etc/pure-ftpd/auth/ -lname "*/PAMAuthentication"  -exec rm -f \{\} \;
find /etc/pure-ftpd/auth/ -lname "*/UnixAuthentication" -exec rm -f \{\} \;

#
# Edit our defaults to set it to "standalone"
#
if [ -f /etc/default/pure-ftpd-common ] ; then
  sed -i -e 's/^STANDALONE_OR_INETD=.*$/STANDALONE_OR_INETD=standalone/' /etc/default/pure-ftpd-common
fi

#
# Disable any inetd service
#
if [ -f /etc/inetd.conf -a -x /usr/sbin/update-inetd ] ; then
  update-inetd --disable ftp
  # Don't really care if this errors.
  invoke-rc.d  openbsd-inetd restart || true
fi

##
# SSL cert generation now in symbiosis-common
##

#
#  Link in the certificate if we have one.
#
if [ ! -e /etc/ssl/private/pure-ftpd.pem ]; then
    ln -s /etc/ssl/ssl.combined /etc/ssl/private/pure-ftpd.pem
fi

if [ -e /etc/ssl/private/pure-ftpd.pem ] ; then
  #
  #  TLS should be enabled if we did the linking.
  #
  echo '1' > /etc/pure-ftpd/conf/TLS

  #
  # Set the TLS cipher suites
  #
  cs="/etc/pure-ftpd/conf/TLSCipherSuite"
  if [ -f "$cs" ] && echo "9939bfb21f5be996e89b8a9e43bab6ba  $cs" | md5sum -c > /dev/null 2>&1 ; then
    # OK this is the Debian ciphers
    echo "I: Updating pure-ftpd cipher suite"
    #
    # If this list changes, remember to change the md5 in the postrm too.
    #
    echo 'ECDHE-RSA-AES128-GCM-SHA256:ECDHE-ECDSA-AES128-GCM-SHA256:ECDHE-RSA-AES256-GCM-SHA384:ECDHE-ECDSA-AES256-GCM-SHA384:DHE-RSA-AES128-GCM-SHA256:DHE-DSS-AES128-GCM-SHA256:kEDH+AESGCM:ECDHE-RSA-AES128-SHA256:ECDHE-ECDSA-AES128-SHA256:ECDHE-RSA-AES128-SHA:ECDHE-ECDSA-AES128-SHA:ECDHE-RSA-AES256-SHA384:ECDHE-ECDSA-AES256-SHA384:ECDHE-RSA-AES256-SHA:ECDHE-ECDSA-AES256-SHA:DHE-RSA-AES128-SHA256:DHE-RSA-AES128-SHA:DHE-DSS-AES128-SHA256:DHE-RSA-AES256-SHA256:DHE-DSS-AES256-SHA:DHE-RSA-AES256-SHA:AES128-GCM-SHA256:AES256-GCM-SHA384:AES128-SHA256:AES256-SHA256:AES128-SHA:AES256-SHA:AES:CAMELLIA:DES-CBC3-SHA:!aNULL:!eNULL:!EXPORT:!DES:!RC4:!MD5:!PSK:!aECDH:!EDH-DSS-DES-CBC3-SHA:!EDH-RSA-DES-CBC3-SHA:!KRB5-DES-CBC3-SHA' > $cs
  fi
fi

#
# Update default runlevels.
#
if [ -x "/etc/init.d/pure-authd" ]; then
  if [ -n "$(which insserv)" ] ; then
    insserv pure-authd
  else
    update-rc.d pure-authd defaults >/dev/null
  fi
fi

#
#  Restart (or start)
#
if [ "$2" = "" ] ; then
  action=start
else
  action=restart
fi

invoke-rc.d pure-authd $action

#
# Restart pure-ftpd, since we've already made config changes.
#

invoke-rc.d pure-ftpd restart 

#DEBHELPER#
exit 0
