# ftpd/debian/postrm
#!/bin/bash
set -e

#
# Remove from the rc levels.
#
if [ "$1" = "purge" ] ; then
  update-rc.d pure-authd remove 
  update-rc.d pure-ftpd remove 
fi

#DEBHELPER#
exit 0
