# ftpd/debian/symbiosis-ftpd.cron.weekly
#!/bin/bash -eu

verbose=""
if [ "$*" = "--verbose" ] || [ "$*" = "-v" ]; then
  verbose="-v "
fi

symbiosis-generate-dhparams $verbose/etc/ssl/private/pure-ftpd-dhparams.pem 
