# ftpd/debian/symbiosis-ftpd.cron.weekly
#!/bin/bash

symbiosis-generate-dhparams /var/lib/symbiosis/ftpd
