# email/cron.d/symbiosis-email
#
# Crontab snippet to scan /srv for new symlinks and writes a config file for
# exim for domain rewrites.
#

*/1 * * * * Debian-exim /usr/sbin/exim_rewrite_scan /etc/exim4/exim4.conf

